# etc/systemd/system/bluetooth.service.d/01-disable-sap.conf
# https://raspberrypi.stackexchange.com/questions/40839/sap-error-on-bluetooth-service-status
[Service]
ExecStart=
ExecStart=/usr/lib/bluetooth/bluetoothd --noplugin=sap
